# gaomon-m10k.service
[Unit]
Description=Pad binding for GAOMON Tablet

[Service]
Type=oneshot
ExecStart=/usr/bin/env gaomon-m10k.sh

[Install]
WantedBy=multi-user.target
